# contrib/applemax-profile/systemd/apex-session.service
[Unit]
Description=AppleMax session bootstrap (DMS primary, Waybar fallback)
PartOf=niri.service
After=niri.service

[Service]
Type=oneshot
ExecStart=%h/.config/apex/scripts/session-up
RemainAfterExit=yes

[Install]
WantedBy=niri.service
